# crawler/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:crawler_main]
command=python main_crawler.py
directory=/workspaces/py_env_research/group05-apache-hive-analyst-socialmedia/crawler
autostart=true
autorestart=true
startretries=3
stderr_logfile=/tmp/crawler_main.err.log
stdout_logfile=/tmp/crawler_main.out.log
user=root
environment=HOME="/root",USER="root"

[program:social_media_scraper]
command=python social_media_scraper.py
directory=/workspaces/py_env_research/group05-apache-hive-analyst-socialmedia/crawler
autostart=false
autorestart=true
startretries=3
stderr_logfile=/tmp/social_scraper.err.log
stdout_logfile=/tmp/social_scraper.out.log
user=root

[group:crawlers]
programs=crawler_main,social_media_scraper
priority=999
